# project_2/server-python.conf
description "Auto-start Python TCP Server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 30

exec /usr/bin/python /home/root/Desktop/server.py